# system/var/lib/opkg/info/aqc-ipa-offload.list
/usr/bin/aqc_gsi_stats
/lib/systemd/system/aqc-ipa-offload.service
