# system/var/lib/opkg/info/avahi-systemd.postinst
set -e
	if [ "x$D" != "x" ]; then
		exit 0
	fi

	systemctl enable avahi-daemon.service
