# system/var/lib/opkg/info/bt-dlkm.list
/lib/systemd/system/bluetooth_power.service
/lib/systemd/
/lib/systemd/system/
/lib/modules/
/lib/modules/4.14.206-perf/extra/
/etc/initscripts/bluetooth_power.sh
/lib/modules/4.14.206-perf/
/lib/modules/4.14.206-perf/extra/bluetooth-power.ko
/lib/systemd/system/multi-user.target.wants/bluetooth_power.service
/lib/systemd/system/multi-user.target.wants/
/etc/initscripts/
