# system/var/lib/opkg/info/data.list
/lib/systemd/system/qcmap_wpa_supplicant@.service
/usr/lib/libqdp.so.0.0.0
/etc/data/mbim_mode
/WEBSERVER/www/QCMAP_Firewall.html
/etc/initscripts/netmgrd
/WEBSERVER/www/QCMAP_Firewall_Help.html
/usr/lib/libqdi.so.0.0.0
/usr/lib/libqcmapipc.so.1.0.0
/lib/systemd/system/port_bridge.service
/etc/data/mobileap_cfg.xml
/etc/initscripts/start_qcmap_wpa_supplicant_le
/WEBSERVER/www/QCMAP_IPV6_Help.html
/WEBSERVER/www/QCMAP_login.html
/etc/udev/rules.d/80-net-setup-link.rules
/lib/systemd/system/multi-user.target.wants/adpld.service
/usr/bin/adpl
/WEBSERVER/
/usr/bin/qmi_ip_multiclient
/lib/systemd/system/QCMAP_Bootup.service
/usr/tests/dsi_netctrl_test
/usr/lib/libnetmgr_common.so.0
/lib/systemd/system/ipacmdiag.service
/etc/tmpfiles.d/qmuxbridge-tmpfilesd.conf
/lib/systemd/system/QCMAP_ConnectionManagerd.service
/lib/systemd/system/qti_pppd.service
/lib/systemd/system/multi-user.target.wants/qmi_ip_multiclientd.service
/lib/systemd/system/qcmap_hostapd@.service
/etc/data/usb/softap_w_dun
/usr/lib/libqcmap_cm.so.1
/etc/tmpfiles.d/atfwd.conf
/WEBSERVER/www/QCMAP_WLAN.html
/etc/data/factory_mobileap_cfg.xml
/usr/bin/netmgrd
/etc/data/l2tp_cfg.xml
/lib/systemd/system/local-fs.target.wants/QCMAP_Bootup.service
/etc/udev/scripts/data_udev_script.sh
/etc/data/factory_mobileap_firewall.xml
/usr/bin/qti_socksv5
/lib/systemd/system/netmgrd.service
/WEBSERVER/www/images/delete.png
/usr/lib/libqcmapipc.so.1
/etc/data/mobileap_cfg.xsd
/usr/bin/radish
/usr/bin/QCMAP_CLI
/lib/systemd/system/qcmap_wpa_cli@.service
/etc/data/wlan_cfg.xml
/etc/data/factory_qti_socksv5_auth.xml
/etc/data/qti_socksv5_auth.xml
/etc/data/hosts
/usr/lib/libqcmaputils.so.1.0.0
/lib/systemd/system/QCMAP_copybind.service
/etc/data/data_iot/carrier_apn_cfg
/WEBSERVER/www/QCMAP_UPNP_DLNA_MDNS.html
/usr/bin/eMBMs_TunnelingModule
/WEBSERVER/www/cgi-bin/qcmap_web_cgi
/lib/systemd/system/multi-user.target.wants/data-init.service
/usr/lib/libqmi_ip.so.1.0.0
/usr/lib/libqcmap_cm.so.1.0.0
/usr/lib/libqdi.so.0
/lib/systemd/system/qcmap_wlan.service
/WEBSERVER/www/QCMAP_Help.html
/WEBSERVER/www/QCMAP_Account_Help.html
/lib/systemd/system/local-fs.target.wants/QCMAP_copybind.service
/lib/systemd/system/adpld.service
/etc/data/factory_qti_socksv5_conf.xml
/WEBSERVER/www/QCMAP_IPV4.html
/etc/initscripts/data-init
/etc/data/qti_socksv5_conf.xml
/usr/lib/libqmi_ip.so.1
/etc/data/netmgr_config.xml
/etc/data/adpl/adpl_config
/etc/data/factory_l2tp_cfg.xml
/WEBSERVER/www/QCMAP_IPV4_Help.html
/usr/lib/libqdp.so.0
/usr/tests/
/etc/data/mobileap_firewall.xml
/WEBSERVER/www/images/
/lib/systemd/system/eMBMs_TunnelingModuled.service
/etc/data/port_bridge/port_bridge.config
/etc/data/factory_wlan_cfg.xml
/usr/bin/ipacmdiag
/usr/bin/cpe_device_manager_demo
/etc/qmi_ip_cfg.xml
/usr/lib/libnetmgr.so.0.0.0
/WEBSERVER/www/QCMAP_WLAN_Help.html
/etc/udev/rules.d/data_udev_rules.rules
/lib/systemd/system/multi-user.target.wants/atfwd.service
/WEBSERVER/www/QCMAP_WWAN.html
/lib/systemd/system/multi-user.target.wants/qmuxbridged.service
/lib/systemd/system/qtid.service
/lib/systemd/system/multi-user.target.wants/ipacmdiag.service
/WEBSERVER/www/QCMAP.html
/usr/lib/libds_atctrl.so.0.0.0
/etc/initscripts/start_port_bridge
/lib/systemd/system/multi-user.target.wants/qtid.service
/etc/udhcpc.d/udhcpc.script
/usr/bin/QCMAP_Web_CLIENT
/etc/data/qmi_ip_cfg.xml
/usr/lib/libdsi_netctrl.so.0
/etc/data/sw_ip_channel_setup_script.sh
/usr/bin/QCMAP_Bootup
/etc/data/dhcp_hosts
/usr/bin/QCMAP_ApInterface
/usr/bin/QCMAP_ConnectionManager
/WEBSERVER/www/QCMAP_WWAN_Help.html
/etc/data/data_iot/bootstrap.ini
/usr/lib/libds_atctrl.so.0
/lib/systemd/system/usb.service.d/
/etc/udev/rules.d/ipacm.rules
/usr/bin/qti
/WEBSERVER/www/cgi-bin/
/usr/lib/libnetmgr_common.so.0.0.0
/WEBSERVER/www/QCMAP_UPNP_DLNA_MDNS_Help.html
/lib/systemd/system/qcmap_load_module@.service
/WEBSERVER/www/QCMAP_IPV6.html
/WEBSERVER/www/images/header.jpg
/usr/bin/QCMAP_StaInterface
/etc/data/data_iot/lwm2m_cfg
/lib/systemd/system/qcmap_hostapd_cli@.service
/WEBSERVER/www/QCMAP_Account.html
/usr/lib/libqcmap_client.so.1
/lib/systemd/system/qmuxbridged.service
/lib/systemd/system/local-fs.target.wants/netmgrd.service
/etc/data/port_bridge/
/lib/systemd/system/ipacm_perf.service
/etc/data/data_iot/
/WEBSERVER/www/QCMAP_NAT.html
/etc/data/mux_mode
/etc/initscripts/start_qcmap_hostapd_le
/WEBSERVER/www/
/lib/systemd/system/qmi_ip_multiclientd.service
/usr/lib/libqcmap_client.so.1.0.0
/usr/bin/qmuxbridge
/WEBSERVER/www/QCMAP_NAT_Help.html
/usr/bin/LwM2M_Client
/usr/bin/port_bridge
/etc/data/adpl/
/lib/systemd/system/atfwd.service
/usr/bin/ssl_demo
/usr/tests/ds_atctrl_test
/lib/systemd/system/data-init.service
/etc/data/usb/
/lib/systemd/system/multi-user.target.wants/port_bridge.service
/usr/lib/libqcmaputils.so.1
/lib/systemd/system/qcmap_wpa_cli_setband@.service
/lib/systemd/system/sockets.target.wants/QCMAP_ConnectionManagerd.service
/usr/bin/qti_ppp
/lib/systemd/system/usb.service.d/usb_data_service.conf
/etc/data/adpl/mhi_enable
/lib/systemd/system/multi-user.target.wants/qti_pppd.service
/etc/initscripts/start_qcmap_wlan_le
/etc/udev/rules.d/99-data-eth.rules
/lib/systemd/system/multi-user.target.wants/ipacm_perf.service
/usr/bin/LwM2M_Client_App
/WEBSERVER/www/cgi-bin/qcmap_auth
/lib/systemd/system/qcmap_wpa_cli_status@.service
/etc/udev/scripts/ipacm_udev.sh
/usr/lib/libdsi_netctrl.so.0.0.0
/WEBSERVER/www/images/update.jpg
/etc/data/data_iot/lwm2m_app_autostart
/usr/lib/libnetmgr.so.0
/etc/data/dsi_config.xml
/etc/tmpfiles.d/data_systemd_tmpfiles.conf
