# system/var/lib/opkg/info/data-ipa-cfg-mgr.list
/lib/systemd/system/local-fs.target.wants/ipacm.service
/usr/lib/libipanat.so.1
/etc/data/ipa/IPACM_cfg.xml
/usr/bin/ipacm
/etc/data/ipa/factory_IPACM_cfg.xml
/etc/data/ipa/
/lib/systemd/system/ipacm.service
/data/misc/
/usr/lib/libipanat.so.1.0.0
/data/misc/ipa/
/etc/tmpfiles.d/ipacm.conf
