# system/var/lib/opkg/info/dbus-1.list
/usr/bin/dbus-send
/usr/share/xml/dbus-1/introspect.dtd
/usr/lib/tmpfiles.d/
/usr/lib/systemd/user/
/usr/lib/systemd/user/sockets.target.wants/dbus.socket
/usr/lib/systemd/user/dbus.socket
/usr/lib/systemd/
/usr/libexec/dbus-daemon-launch-helper
/usr/share/dbus-1/session.conf
/usr/share/dbus-1/system.conf
/usr/share/xml/
/usr/bin/dbus-uuidgen
/usr/lib/sysusers.d/dbus.conf
/lib/systemd/system/dbus.socket
/lib/systemd/system/sockets.target.wants/dbus.socket
/usr/bin/dbus-daemon
/usr/share/dbus-1/system.d/extra-users-reboot.conf
/lib/systemd/system/dbus.target.wants/dbus.socket
/usr/share/xml/dbus-1/
/usr/share/xml/dbus-1/busconfig.dtd
/etc/dbus-1/
/usr/bin/dbus-cleanup-sockets
/usr/share/dbus-1/system.d/
/var/lib/dbus/
/lib/systemd/system/sockets.target.wants/dbus.service
/usr/share/dbus-1/session.d/
/usr/share/dbus-1/
/usr/lib/sysusers.d/
/usr/bin/dbus-run-session
/lib/systemd/system/dbus.service.d/
/usr/bin/dbus-launch
/usr/share/dbus-1/system-services/
/usr/bin/dbus-update-activation-environment
/lib/systemd/system/dbus.service
/usr/lib/systemd/user/dbus.service
/etc/default/volatiles/99_dbus
/usr/share/dbus-1/system.d/extra-users.conf
/usr/share/dbus-1/services/
/usr/lib/systemd/user/sockets.target.wants/
/etc/dbus-1/session.conf
/etc/dbus-1/system.conf
/lib/systemd/system/dbus.target.wants/
/lib/systemd/system/dbus.service.d/dbus.conf
/usr/bin/dbus-monitor
/usr/lib/tmpfiles.d/dbus.conf
