# system/var/lib/opkg/info/ddclient.list
/etc/data/ddclient.conf
/lib/systemd/system/ddclientd.service
/usr/sbin/ddclient
