# system/var/lib/opkg/info/emac-dwc-eqos.postinst
set -e
   if false; then
    [ -n "$D" ] && OPT="-r $D" || OPT="-s"
        update-rc.d $OPT -f emac_dwc_eqos_start_stop_le remove
        update-rc.d $OPT emac_dwc_eqos_start_stop_le start 37 S . stop 63 0 1 6 .
        update-rc.d $OPT -f setup_avtp_routing_le remove
        update-rc.d $OPT setup_avtp_routing_le start 91 S . stop 9 0 1 6 .
    fi
