# system/var/lib/opkg/info/embms-kernel.list
/lib/systemd/system/embmsd.service
/etc/initscripts/start_embms_le
/lib/systemd/system/multi-user.target.wants/embmsd.service
