# system/var/lib/opkg/info/init-mss.list
/etc/tmpfiles.d/
/sbin/init_mss
/etc/tmpfiles.d/init_mss.conf
/lib/systemd/system/init_sys_mss.service
/etc/udev/rules.d/init_mss.rules
/lib/systemd/system/sockets.target.wants/init_sys_mss.service
