# system/var/lib/opkg/info/libdiag1.postinst
set -e
       if false; then
        [ -n "$D" ] && OPT="-r $D" || OPT="-s"

        update-rc.d $OPT -f chgrp-diag remove
        update-rc.d $OPT chgrp-diag start 15 2 3 4 5 .
       fi
if [ x"$D" = "x" ]; then
	if [ -x /sbin/ldconfig ]; then /sbin/ldconfig ; fi
fi
