# system/var/lib/opkg/info/mcm-core.list
/usr/bin/MCM_MOBILEAP_ConnectionManager
/usr/bin/mcm_data_srv
/usr/lib/libmcm_ssr_util.so.0
/usr/lib/libmcm_ssr_util.so.0.0.0
/usr/bin/MCM_MOBILEAP_CLI
/usr/lib/libmcm.so.0.0.0
/usr/bin/MCM_atcop_svc
/usr/lib/libmcm_log_util.so.0
/usr/bin/uim_test_client
/lib/systemd/system/mcm-core.service
/etc/initscripts/start_mcm_data_srv_le
/etc/initscripts/start_MCM_MOBILEAP_ConnectionManager_le
/usr/bin/mcm_ril_service
/usr/lib/libmcmipc.so.0
/etc/tmpfiles.d/mcm-core_systemd_tmpfiles.conf
/usr/tests/mcm_data_test
/usr/lib/libmcmipc.so.0.0.0
/usr/lib/libmcm.so.0
/usr/bin/MCM_ATCOP_CLI
/usr/bin/IoEConsoleClient
/usr/lib/libmcm_log_util.so.0.0.0
/etc/initscripts/start_MCM_atcop_svc_le
