# system/var/lib/opkg/info/miniupnpd.list
/usr/sbin/miniupnpd
/etc/data/miniupnpd/
/lib/systemd/system/miniupnpd.service
/etc/data/miniupnpd/miniupnpd.conf
