# system/var/lib/opkg/info/modem-shutdown.list
/etc/initscripts/modem-shutdown
/lib/systemd/system/poweroff.target.wants/modem_shutdown.service
/bin/modem-shutdown
/lib/systemd/system/poweroff.target.wants/
/lib/systemd/system/reboot.target.wants/modem_shutdown.service
/lib/systemd/system/modem_shutdown.service
/lib/systemd/system/reboot.target.wants/
