# system/var/lib/opkg/info/pimd.list
/etc/pimd.conf
/usr/sbin/pimd
/lib/systemd/system/pimd.service
