# system/var/lib/opkg/info/ppp.list
/usr/sbin/pppd
/etc/ppp/ipv6-up
/usr/lib/pppd/2.4.7/passwordfd.so
/etc/chatscripts/pap
/etc/ppp/pap-secrets
/etc/ppp/ip-down.d/
/etc/ppp/ip-up.d/08setupdns
/usr/lib/pppd/2.4.7/winbind.so
/usr/lib/pppd/2.4.7/radattr.so
/usr/lib/pppd/2.4.7/pppoatm.so
/usr/lib/pppd/
/etc/ppp/
/usr/lib/pppd/2.4.7/radius.so
/usr/lib/pppd/2.4.7/radrealms.so
/lib/systemd/system/pppd_service.service
/etc/ppp/ip-down.d/92removedns
/etc/ppp/ip-down
/usr/bin/pon
/usr/lib/pppd/2.4.7/minconn.so
/usr/lib/pppd/2.4.7/passprompt.so
/etc/ppp/ppp_on_boot
/etc/ppp/ip-up
/etc/ppp/ip-up.d/
/usr/lib/pppd/2.4.7/openl2tp.so
/etc/ppp/options
/etc/disconnect
/etc/chatscripts/
/usr/lib/pppd/2.4.7/
/etc/ppp/chap-secrets
/usr/sbin/chat
/usr/lib/pppd/2.4.7/rp-pppoe.so
/usr/bin/poff
/etc/ppp/peers
/lib/systemd/system/ppp@.service
/usr/lib/pppd/2.4.7/pppol2tp.so
