# system/var/lib/opkg/info/qce.list
/etc/systemd/system/multi-user.target.wants/
/etc/systemd/
/etc/systemd/system/qce_dlkm.service
/etc/systemd/system/
/etc/initscripts/qce_dlkm
/etc/systemd/system/multi-user.target.wants/qce_dlkm.service
