# system/var/lib/opkg/info/ssreq-server.list
/usr/bin/ssreq_server
/lib/systemd/system/ssreq-server.service
/lib/systemd/system/local-fs.target.wants/ssreq-server.service
