# system/var/lib/opkg/info/start-scripts-sfs-config.list
/etc/systemd/system/multi-user.target.wants/sfsconfig.service
/sbin/sfs_config
/etc/systemd/system/sfsconfig.service
