# system/var/lib/opkg/info/start-scripts-qseecomd-daemon.list
/etc/udev/
/sbin/qseecomd
/etc/udev/rules.d/99-qseecomd.rules
/lib/systemd/system/qseecomd.service
/etc/udev/rules.d/
/lib/systemd/system/local-fs.target.requires/
/lib/systemd/system/local-fs.target.requires/qseecomd.service
