# system/var/lib/opkg/info/start-scripts-misc-daemon.list
/etc/systemd/system/multi-user.target.wants/misc_daemon.service
/etc/initscripts/misc-daemon
/etc/systemd/system/misc_daemon.service
