# system/var/lib/opkg/info/subsystem-ramdump.list
/lib/systemd/system/subsystem-ramdump.service
/usr/bin/subsystem_ramdump
