# system/var/lib/opkg/info/systemd-machine-units.list
/lib/systemd/system/systemrw.mount
/lib/systemd/system/local-fs.target.requires/systemrw.mount
/lib/systemd/system/var-volatile.mount
/lib/systemd/system/local-fs.target.wants/var-volatile.mount
/etc/initscripts/firmware-ubi-mount.sh
/lib/systemd/system/persist-ubi.mount
/lib/systemd/system/systemrw.mount.d/
/lib/systemd/system/multi-user.target.wants/label-persist.service
/lib/systemd/system/local-fs.target.wants/
/lib/systemd/system/systemrw.mount.d/systemrw.conf
/lib/systemd/system/multi-user.target.wants/cache.mount
/lib/systemd/system/data-ubi.mount
/lib/systemd/system/local-fs.target.requires/firmware-mount.service
/lib/systemd/system/sysinit.target.wants/persist.mount
/lib/systemd/system/cache-ubi.mount
/lib/systemd/system/sysinit.target.wants/
/lib/systemd/system/firmware-mount.service
/lib/systemd/system/persist.mount
/lib/systemd/system/multi-user.target.wants/label-cache.service
/lib/systemd/system/fw-ubi-mount.service
/lib/systemd/system/firmware.mount
/lib/systemd/system/local-fs.target.requires/firmware.mount
/lib/systemd/system/systemrw-ubi.mount
/lib/systemd/system/cache.mount
/lib/systemd/system/data.mount
/lib/systemd/system/label-data.service
/lib/systemd/system/label-cache.service
/lib/systemd/system/local-fs.target.wants/data.mount
/lib/systemd/system/label-persist.service
/lib/systemd/system/local-fs.target.wants/label-data.service
