# system/var/lib/opkg/info/time-services.list
/etc/initscripts/time_serviced
/etc/systemd/system/time_serviced.service.d/
/etc/systemd/system/time_serviced.service.d/time-services.conf
/etc/udev/rules.d/time-services.rules
/usr/bin/time_daemon
/etc/systemd/system/multi-user.target.wants/time_serviced.service
/etc/systemd/system/time_serviced.service
/etc/tmpfiles.d/time-services-tmpfilesd.conf
