# system/var/lib/opkg/info/volatile-binds.list
/sbin/umount-copybind
/lib/systemd/system/systemrw-rt_tables.service
/lib/systemd/system/systemrw-boot_hsic_comp.service
/sbin/mount-copybind
/lib/systemd/system/systemrw-allplay.service
/lib/systemd/system/systemrw-build.prop.service
/lib/systemd/system/systemrw-adb_devid.service
/lib/systemd/system/systemrw-boot_hsusb_comp.service
/lib/systemd/system/systemrw-data-usb.service
/lib/systemd/system/systemrw-data-adpl.service
/lib/systemd/system/systemrw-data-ipa.service
/lib/systemd/system/systemrw-misc-wifi.service
/lib/systemd/system/systemrw-bluetooth.service
/lib/systemd/system/systemrw-data.service
/lib/systemd/system/systemrw-data-miniupnpd.service
