# system/var/lib/opkg/info/wlan-conf.list
/etc/cld-hostapd-wlan1.sh
/etc/misc/wifi/hostapd-wlan2.conf
/etc/misc/wifi/
/etc/cld-hostapd-wlan2.sh
/etc/misc/
/etc/init.d/start_cnss_le
/etc/sysctl.conf
/etc/misc/wifi/hostapd.conf
/etc/misc/wifi/wpa_supplicant.conf
/etc/misc/wifi/sta_mode_hostapd.conf
/lib/firmware/wlan/qca_cld/WCNSS_qcom_cfg.ini
/etc/misc/wifi/WCNSS_qcom_cfg.ini
/etc/cld-hostapd-wlan0.sh
/etc/misc/wifi/hostapd-wlan1.conf
